# supervisor/mysql.conf
[program:mysql]
command = /usr/sbin/mysql
autostart = true